# backup-sql.cron
0 0 * * * mysqldump -u root -pchangeit txtbudget > /disk1/sqlbackup/txtbudget_$(date +\%Y-\%m-\%d).sql
